# scripts/com.chg.imessage-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.chg.imessage-daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>python3</string>
        <string>/Users/jasricha/Documents/GitHub/chief_of_staff/scripts/imessage-daemon.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/jasricha</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>IMESSAGE_DAEMON_ENABLED</key>
        <string>true</string>
        <key>IMESSAGE_DAEMON_REPLY_HANDLE</key>
        <string></string>
        <key>IMESSAGE_DAEMON_MONITORED_CONVERSATION</key>
        <string></string>
        <key>ANTHROPIC_API_KEY</key>
        <string></string>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/jasricha/Library/Logs/Jarvis/imessage-daemon.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jasricha/Library/Logs/Jarvis/imessage-daemon.err</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
</dict>
</plist>
